# dbibackend/usb-monitor.service
[Unit]
Description=Switch USB Monitor — auto-start dbibackend containers
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/usb_monitor.sh
Restart=on-failure
RestartSec=5
Environment=GAMES_DIR=/volume1/Switch
Environment=QUEUE_DIR=/tmp
Environment=SERVER_URL=http://localhost:8080
StandardOutput=journal
StandardError=journal
# Pastikan cleanup berjalan saat service stop
KillMode=process
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
